# .config/systemd/user/console.target
[Unit]
Description=Console target

[Install]
Alias=default.target
