# .config/systemd/user/graphical.target
[Unit]
Description=Graphical target - graphical daemons, xinitrc stuff etc.
After=xorg.target
