# .config/systemd/user/graphical.target
[Unit]
Description=Xinitrc Stuff
After=xorg.target
BindsTo=xorg.target
