# .config/systemd/user/mpd-notify.service
[Unit]
Description=MPD notifications using libnotify
After=mpd.service
Requires=mpd.service
Requires=dbus.socket

[Service]
Environment=DISPLAY=:0
ExecStart=/home/lahwaacz/bin/mpd-notify
Restart=always

[Install]
WantedBy=graphical.target
