# .config/systemd/user/timer-daily.target
[Unit]
Description=Daily Timer Target
StopWhenUnneeded=yes
