# .config/systemd/user/timer-weekly.target
[Unit]
Description=Weekly Timer Target
StopWhenUnneeded=yes
