# .config/systemd/user/tinyterm.service
[Unit]
Description=qpdfview

[Service]
UMask=0077
Environment=DISPLAY=:0
ExecStart=/usr/bin/tinyterm

[Install]
WantedBy=graphical.target
