# .config/systemd/user/urxvtd.service
[Unit]
Description=Urxvt Daemon

[Service]
UMask=0077
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=graphical.target
