# .config/systemd/user/timer-weekly.target.wants/backup-kmlinux.service
[Unit]
Description=Backup kmlinux files

[Service]
Type=oneshot
ExecStart=/home/lahwaacz/bin/backup-wordpress-kmlinux.sh --quiet
UMask=0077
