# .config/systemd/user/timer-weekly.target.wants/rmshit.service
[Unit]
Description=rmshit

[Service]
Type=oneshot
ExecStart=/usr/bin/vendor_perl/rmshit -c
UMask=0077
