# .config/systemd/user/old-services/dropbox.timer
[Unit]
Description=Dropbox Timer

[Timer]
OnActiveSec=5min
Unit=dropbox.service

[Install]
WantedBy=timers.target
