# .config/systemd/user/old-services/wm.target
[Unit]
Description=Window manager target
Wants=graphical.target
Before=graphical.target
Requires=dbus.socket
AllowIsolate=true
BindsTo=i3.service

[Install]
Alias=default.target
