# .config/systemd/user/old-services/wmfs-statusbar.service
[Unit]
Description=WMFS2 statusbar
After=wmfs.service
Requisite=wmfs.service

[Service]
UMask=0077
Environment=DISPLAY=:0
ExecStart=/home/lahwaacz/bin/wmfs-statusbar -c %h/.config/wmfs/statusbar.conf
Restart=always

[Install]
WantedBy=graphical.target
