# .config/systemd/user/timer-daily.target.wants/backup-kmlinux.service
[Unit]
Description=Backup kmlinux files

[Service]
Type=oneshot
ExecStart=/home/lahwaacz/bin/backup-wordpress-kmlinux.sh -q
UMask=0077
